# ops/apollo-browser-worker.service
[Unit]
Description=Apollo managed browser worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=apollo-agent
WorkingDirectory=/opt/apollo-web-agent/current/browser-worker
EnvironmentFile=/opt/apollo-web-agent/shared/.env
Environment=HOME=/opt/apollo-web-agent/shared/browser-worker-home
ExecStart=/opt/apollo-web-agent/shared/browser-worker-venv/bin/python server.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/apollo-web-agent/shared
MemoryHigh=1500M
MemoryMax=2G
TasksMax=256
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
