# ops/web-agent-frpc.service
[Unit]
Description=Apollo Web Agent FRP
After=network-online.target frps.service nginx.service web-agent.service
Wants=network-online.target

[Service]
Type=simple
User=apollo-agent
Group=apollo-agent
ExecStart=/opt/apollo-web-agent/bin/frpc -c /opt/apollo-web-agent/shared/frp/frpc.toml
Restart=always
RestartSec=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/apollo-web-agent/bin
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
