# ops/web-agent.service
[Unit]
Description=Apollo Web Agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=apollo-agent
Group=apollo-agent
WorkingDirectory=/opt/apollo-web-agent/current
EnvironmentFile=/opt/apollo-web-agent/shared/.env
Environment=PORT=9130
Environment=HOME=/opt/apollo-web-agent/shared/home
Environment=SHELL=/bin/bash
ExecStart=/usr/bin/node server-dist/server/production.js
Restart=always
RestartSec=3
TimeoutStopSec=25
KillSignal=SIGTERM
OOMPolicy=stop
MemoryHigh=2G
MemoryMax=2500M
TasksMax=512
LimitNOFILE=65536
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/apollo-web-agent/shared
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
