# ops/web-agent-backup.timer
[Unit]
Description=Daily Apollo Web Agent data backup

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=20m
Persistent=true
Unit=web-agent-backup.service

[Install]
WantedBy=timers.target
