# ops/web-agent-backup.service
[Unit]
Description=Apollo Web Agent data backup
After=web-agent.service

[Service]
Type=oneshot
User=apollo-agent
Group=apollo-agent
WorkingDirectory=/opt/apollo-web-agent/current
EnvironmentFile=/opt/apollo-web-agent/shared/.env
ExecStart=/usr/bin/node scripts/backup-production.mjs --source /opt/apollo-web-agent/shared --destination /opt/apollo-web-agent/backups
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/apollo-web-agent/shared /opt/apollo-web-agent/backups
